# templates/default/etc/init/pgbouncer.conf.erb
description "PgBouncer PostgreSQL connection pool for db alias <%= @name %>"
start on runlevel [2345]
stop on runlevel [!2345]

env USER=<%= @user %>
env RUNDIR=<%= @instance['pid_dir'] %>
env DAEMON="/usr/sbin/pgbouncer"
env DAEMON_OPTS="-d -R /etc/pgbouncer/pgbouncer-<%= @name %>.ini"

pre-start script
    mkdir -p $RUNDIR || true
    chown $USER:$USER $RUNDIR || true
end script

limit nofile <%= @instance['soft_limit'] %> <%= @instance['hard_limit'] %>
expect daemon
respawn

exec start-stop-daemon --start --quiet --pidfile $RUNDIR/$JOB.pid --chuid $USER \
		--exec $DAEMON -- $DAEMON_OPTS || true
